# kiwi/srv-kiwi.automount
Description=Automount /srv/kiwi on demand

[Automount]
Where=/srv/kiwi

[Install]
WantedBy=multi-user.target
